# systemd/user/anydesk.service
[Unit]
Description=Anydesk
After=network-online.target

[Service]
ExecStart=/usr/sbin/anydesk --control

[Install]
WantedBy=default.target
