# systemd/user/anydesk.socket
[Unit]
Description=Anydesk socket

[Socket]
ListenStream=0.0.0.0:7070
Service=anydesk.service

[Install]
WantedBy=sockets.target
