# systemd/user/certsrenew.timer
[Unit]
Description=Renew SSL certs

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
