# systemd/user/mpDris2.socket
[Unit]
Description=MPRIS2 support for MPD

[Socket]
ListenStream=6600
Service=mpDris2.service
